# configs/etc/systemd/system/ufw.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/iptables-clear
